# pi/systemd/minstrel-server.service
[Unit]
Description=Minstrel Codex Static File Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/minstrelcodex
ExecStart=/usr/bin/env npx serve apps/web/dist -l 3000 -s
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
